# core/plists/com.bob.vocalkord.scan.plist
#!/bin/bash
# ∴ plist for PARSE+FORGE_VOCALKORD.sh — Sanity check for blessing audio

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bob.vocalkord.scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
     <string>/Users/Users/G/blurOS/_sense/PARSE+FORGE_VOCALKORD.sh</string>
  <string>--full</string>
</array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key> <integer>6</integer> <!-- Friday -->
    <key>Hour</key> <integer>11</integer>
    <key>Minute</key> <integer>11</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/G/blur/TEHE/bob_kords_warmed.txt</string>
  <key>StandardErrorPath</key>
  <string>/Users/G/blur/MEEP/bob_kords_cut.err</string>


  <key>RunAtLoad</key> <false/>
</dict>
</plist>
